# ubuntu-intune/system_files/etc/systemd/user/microsoft-identity-broker.service.d/override.conf
[Service]
EnvironmentFile=-/etc/intuneme-host-display
